# packages/webOS/luna-service2/files/systemd/ls-hubd_private.service
[Unit]
Description=Private luna-service2 bus

[Service]
ExecStart=/usr/sbin/ls-hubd --conf /etc/luna-service2/ls-private.conf --pmloglib

[Install]
#WantedBy=multi-user.target
