# packages/webOS/luna-service2/files/systemd/ls-hubd_public.service
[Unit]
Description=Public luna-service2 bus
After=ls-hubd_private.service

[Service]
ExecStart=/usr/sbin/ls-hubd --public --conf /etc/luna-service2/ls-public.conf --pmloglib

[Install]
#WantedBy=multi-user.target
